# hack/node-agent/node-agent.service
[Unit]
Description=Web Command Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/srv/node-agent
EnvironmentFile=-/srv/node-agent/agent.env
Environment="TMPDIR=/srv/node-agent"
ExecStart=/srv/node-agent/node-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target